# etc/periodic/daily/Makefile
# $FreeBSD: src/etc/periodic/daily/Makefile,v 1.6 2001/02/19 07:12:37 ache Exp $

BIN=	340.uucp \
	410.status-uucp

.include <bsd.prog.mk>
